# packaging/aznfs/DEBIAN/postinst
#!/bin/bash

# Exit on error.
set -e

chmod +x /sbin/aznfswatchdog
chmod +x /sbin/mount.aznfs

systemctl daemon-reload
systemctl enable aznfswatchdog
systemctl start aznfswatchdog
systemctl status aznfswatchdog
