# packaging/aznfs/DEBIAN/postinst
#!/bin/bash

# --------------------------------------------------------------------------------------------
# Copyright (c) Microsoft Corporation. All rights reserved.
# Licensed under the MIT License. See License.txt in the project root for license information.
# --------------------------------------------------------------------------------------------

export NCURSES_NO_UTF8_ACS=1

# Exit on error.
set -e

FLAG_FILE="/tmp/.update_in_progress_from_watchdog.flag"
CONFIG_FILE="/opt/microsoft/aznfs/data/config"
AUTO_UPDATE_AZNFS="false"

parse_user_config()
{
    if [ ! -f "$CONFIG_FILE" ]; then
        echo "$CONFIG_FILE not found. Please make sure it is present."
    fi

    # Read the value of AUTO_UPDATE_AZNFS from the configuration file and convert to lowercase for easy comparison later.
    AUTO_UPDATE_AZNFS=$(egrep -o '^AUTO_UPDATE_AZNFS[[:space:]]*=[[:space:]]*[^[:space:]]*' "$CONFIG_FILE" | tr -d '[:blank:]' | cut -d '=' -f2)
    AUTO_UPDATE_AZNFS=${AUTO_UPDATE_AZNFS,,}
}

user_consent_for_auto_update()
{
    # Check the value of AUTO_UPDATE_AZNFS in the config file
    parse_user_config

    if [ "$AUTO_UPDATE_AZNFS" == "true" ]; then
        return 0
    fi
#*******************************************************************************************************
#     # Get the terminal window dimensions
#     rows=$(tput lines)
#     cols=$(tput cols)

#     # To Keep dialog box size based on screen dimensions
#     height=$((rows * 30 / 100))
#     width=$((cols * 60 / 100))

#     title="Auto-Update Configuration"
#     auto_update_prompt="Do you wish to enable automatic updates for AZNFS to ensure you stay up-to-date with the \
#                         latest features, improvements, and security patches? We recommend enabling automatic updates \
#                         to ensure you have the best AZNFS experience. If you choose to enable automatic updates, \
#                         the AZNFS will periodically check for updates and apply them automatically."

#     sed -i '/AUTO_UPDATE_AZNFS/d' "$CONFIG_FILE"

#     # Run the dialog command and capture its exit code.
#     dialog --timeout 300 --default-button yes --title "$title" --yesno "$auto_update_prompt" $height $width
#     dialog_exit_code=$?

#     if [ $dialog_exit_code -eq 0 ] || [ $dialog_exit_code -eq 255 ]; then
#         # User chose "Yes" or it timed out (consider it as "true")
#         echo "AUTO_UPDATE_AZNFS=true" > "$CONFIG_FILE"
#     else
#         # User chose "No"
#         echo "AUTO_UPDATE_AZNFS=false" > "$CONFIG_FILE"
#     fi
#*******************************************************************************************************

# Check if dialog is installed
if ! command -v dialog > /dev/null; then
    echo "Error: 'dialog' is not installed. Please install it and try again."
    exit 1
fi

# Debugging
echo "Before dialog"
dialog --msgbox "I'm a dialog box message" 10 30 < /dev/tty
dialog_exit_code=$?
echo "After dialog"

# Debugging: Output the exit code
echo "Exit code: $dialog_exit_code"
}

# Set appropriate permissions.
chmod 0755 /opt/microsoft/aznfs/
chmod 0755 /usr/sbin/aznfswatchdog
chmod 0755 /opt/microsoft/aznfs/mountscript.sh
chmod 0755 /opt/microsoft/aznfs/aznfs_install.sh
chmod 0644 /opt/microsoft/aznfs/common.sh

# Set suid bit for mount.aznfs to allow mount for non-super user.
chmod 4755 /sbin/mount.aznfs

# Create data directory for holding mountmap and log file. 
mkdir -p /opt/microsoft/aznfs/data
chmod 0755 /opt/microsoft/aznfs/data

# Move the mountmap, aznfs.log and randbytes files to new path in case these files exists and package is being upgraded.
if [ -f /opt/microsoft/aznfs/mountmap ]; then
        chattr -f -i /opt/microsoft/aznfs/mountmap
        mv -vf /opt/microsoft/aznfs/mountmap /opt/microsoft/aznfs/data/
        chattr -f +i /opt/microsoft/aznfs/data/mountmap
fi

if [ -f /opt/microsoft/aznfs/aznfs.log ]; then
        mv -vf /opt/microsoft/aznfs/aznfs.log /opt/microsoft/aznfs/data/
fi

if [ -f /opt/microsoft/aznfs/randbytes ]; then
        chattr -f -i /opt/microsoft/aznfs/randbytes
        mv -vf /opt/microsoft/aznfs/randbytes /opt/microsoft/aznfs/data/
        chattr -f +i /opt/microsoft/aznfs/data/randbytes
fi

# Check if the config file exists; if not, create it.
if [ ! -f /opt/microsoft/aznfs/data/config ]; then
        # Create the config file and set default AUTO_UPDATE_AZNFS=false inside it.
        echo "AUTO_UPDATE_AZNFS=false" > /opt/microsoft/aznfs/data/config

        # Set the permissions for the config file.
        chmod 0644 /opt/microsoft/aznfs/data/config
fi

# If it's an auto update triggered by aznfswatchdog, don't restart watchdog.
if [ ! -f "$FLAG_FILE" ]; then

        user_consent_for_auto_update

        systemctl daemon-reload
        systemctl enable aznfswatchdog
        systemctl start aznfswatchdog
else
        # Clean up the update in progress flag file.
        rm -f "$FLAG_FILE"
fi