# packaging/aznfs/DEBIAN/prerm
#!/bin/bash

# --------------------------------------------------------------------------------------------
# Copyright (c) Microsoft Corporation. All rights reserved.
# Licensed under the MIT License. See License.txt in the project root for license information.
# --------------------------------------------------------------------------------------------

# Exit on error.
set -e

if [ "$1" == 'remove' -o "$1" == 'purge' ]; then
	# Verify if any existing mounts are there, warn the user about this.
	existing_mounts=$(cat /opt/microsoft/aznfs/mountmap 2>/dev/null | egrep '^\S+' | wc -l)
	if [ $existing_mounts -ne 0 ]; then
		echo "There are existing Azure Blob NFS mounts using aznfs mount helper, they will not be tracked!"
		read -n 1 -p "Are you sure you want to continue? [Y/n] " result
		echo
		if [ -n "$result" -a "$result" != "y" -a "$result" != "Y" ]; then
			echo "Removal aborted!"
			exit 0
		fi
	fi
fi

# Stop aznfswatchdog service and clean it.
systemctl stop aznfswatchdog
systemctl disable aznfswatchdog
echo "Stopped aznfswatchdog service."