# packaging/aznfs/DEBIAN/prerm
#!/bin/bash

# Exit on error.
set -e

systemctl stop aznfswatchdog
systemctl disable aznfswatchdog
systemctl daemon-reload
echo "Stopped aznfswatchdog service."